# dnsmasq-blocklist.timer
[Unit]
Description=Daily timer for dnsmasq-blocklist service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

